# upstart/testjob.conf
description "A test job file for experimenting with Upstart"
author "Oğuz Yarımtepe"

start on [2345]
exec echo Test Job ran at  `date` >> /var/log/testjob.log

